# com.blogagent.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blogagent.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/choo/.gemini/antigravity/scratch/venv/bin/python</string>
        <string>/Users/choo/.gemini/antigravity/scratch/blog_agent/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/choo/.gemini/antigravity/scratch/blog_agent</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>
    <key>StandardOutPath</key>
    <string>/Users/choo/.gemini/antigravity/scratch/blog_agent/scheduler_run.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/choo/.gemini/antigravity/scratch/blog_agent/scheduler_run.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
